# launchd/ai.lictor.monthly-metrics.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.lictor.monthly-metrics</string>

  <!-- Wraps the metrics-generator script. Loads ~/.lictor/secrets.env so API
       keys are available; runs in the user's full launchd context. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/raffa/Lictor/scripts/run-monthly-metrics.sh</string>
  </array>

  <!-- 1st of every month at 06:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/lictor-monthly-metrics-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/lictor-monthly-metrics-launchd.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
