# launchd/ai.lictor.il-patrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.lictor.il-patrol</string>

  <!-- Israel first-priority daily VALIDATED scan. Refreshes the IL master,
       runs the verifier-gated (zero-FP) scan, then refreshes the Globe warm
       map for IL. Draft-only — writes a local ledger, sends nothing. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/raffa/Lictor/scripts/il-daily-patrol.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/raffa/Lictor</string>

  <!-- Daily at 03:00 local — its own dedicated pass, ahead of the global
       3×-daily globe (02/10/18) so Israel always gets a clean full sweep. -->
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>

  <key>StandardOutPath</key>
  <string>/tmp/lictor-il-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/lictor-il-daily.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
